# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-scheduler/finish
#!/command/with-contenv /usr/bin/bash
# shellcheck shell=bash

if [[ $1 -eq 0 ]]; then
	echo "[svc-scheduler] Service stopped cleanly (exit code: $1, signal: $2)"
else
	echo "[svc-scheduler] Service exited with error (exit code: $1, signal: $2)" >&2
fi
